# .local/share/systemd/user/speedtest.service
[Unit]
Description=Measuring download speed
RefuseManualStart=no
RefuseManualStop=yes
 
[Service]
Type=oneshot
ExecStart=/usr/local/bin/speedtest.sh
